# etc/supervisord.d/mariadb.ini
[program:mariadb]
command = /usr/bin/mysqld_safe --datadir='/var/lib/mysql'
user=mysql
group=mysql
autostart=false
autorestart=true
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0